# container/services/etcd/run
#!/bin/sh
exec 2>&1

CONFIG="/container/config/etcd/etcd.yml"
EXTRAOPTS="--logger=zap"

exec chpst -u etcd /usr/bin/etcd --config-file=$CONFIG $EXTRAOPTS 2>&1 | logger --tag etcd
