# container/services/etcd/run
#!/bin/sh
exec 2>&1

CONFIG="/container/config/etcd/etcd.config"
EXTRAOPTS=""

exec chpst -u etcd /usr/bin/etcd --config $CONFIG $EXTRAOPTS
